# deploy/masignals-api.service
[Unit]
Description=MA-Signals API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/financial-news-pipeline
EnvironmentFile=/opt/financial-news-pipeline/.env
# 127.0.0.1 par defaut : l'API n'est pas exposee a l'exterieur de la VM.
# Pour un acces distant : tunnel SSH (ssh -L 8000:localhost:8000 vm) ou reverse
# proxy TLS, et definir API_KEY dans le .env avant de repasser en 0.0.0.0.
ExecStart=/opt/financial-news-pipeline/.venv/bin/uvicorn ma_signals.api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10
User=masignals

[Install]
WantedBy=multi-user.target
